# bot/Config/supervisord.conf
[supervisord]
nodaemon=true
silent=true
logfile=../data/logs/supervisord.log
loglevel=debug
pidfile=../data/supervisord.pid

[program:api]
command=gunicorn -c ./Config/gunicorn.conf.py API:app

[program:polybot]
command=python3 -m PolyBot
